# scripts/optimization/tunix-optimize.service
[Unit]
Description=TUNIX System Optimization Service
After=network.target tunix-power.service
Requires=tunix-power.service tunix-system-control.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/lib/tunix/hardware-profile.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target