# scripts/optimization/tunix-network.service
[Unit]
Description=TUNIX Network Optimization Service
After=network-online.target
Wants=network-online.target
Requires=tunix-system-control.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/tunix/network_optimizer.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target